# docsrc/de/modules/installation/attachments/podman/quadlet/titracka-redis.container
[Unit]
Description=Redis for Cocard (production)

[Container]
Image=docker.io/library/redis:latest
ContainerName=titracka-redis
Volume=titracka-redis.volume:/data
Network=titracka.network
Label=traefik.enable=false

[Service]
Restart=always
# pulling an image may exceed the default timeout
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target default.target
